# os/config/update-motd.d/10-help-text
#!/bin/sh
# Koheron: suppress Ubuntu help text in the login banner.
exit 0
