# os/config/update-motd.d/50-motd-news
#!/bin/sh
# Koheron: disable Ubuntu motd news snippets.
exit 0
