# os/config/update-motd.d/60-unminimize
#!/bin/sh
# Disable the Ubuntu "unminimize" prompt in the MOTD
exit 0
