# os/systemd/koheron-server.service
[Unit]
Description=Koheron TCP/Websocket server
After=network.target

[Service]
Type=forking
ExecStartPre=mkdir -p /tmp/instruments && /bin/cp -n /usr/local/koheron-server/kserverd /tmp/instrument/kserverd
ExecStart=/tmp/instruments/kserverd -c /usr/local/koheron-server/koheron-server.conf
ExecStop=/usr/bin/pkill -SIGINT kserverd
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target

