# os/systemd/koheron-server.service
[Unit]
Description=Koheron TCP/Websocket server
DefaultDependencies=no
# minimal boot deps to be safe; no network wait
After=sysinit.target local-fs.target tmp.mount unzip-default-instrument.service
Before=basic.target
Requires=unzip-default-instrument.service

[Service]
Type=notify
NotifyAccess=all
ExecStart=/tmp/live-instrument/serverd
ExecStop=/usr/bin/pkill -SIGINT serverd
KillSignal=SIGKILL
StartLimitInterval=0

[Install]
WantedBy=basic.target