# config/tcp-server.service
[Unit]
Description=Koheron TCP server
After=network.target

[Service]
ExecStart=/usr/local/tcp-server/kserverd -c /usr/local/tcp-server/kserver.conf
Type=notify

[Install]
WantedBy=multi-user.target
